# debian/prerm
#!/bin/bash
set -e
service my-bindfs-mounts stop
update-rc.d my-bindfs-mounts remove
